# workspace/crontab/laradock
* * * * * laradock php /var/www/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock php /var/www/ampv2/artisan schedule:run >> /dev/null 2>&1
